# thinkled.service
[Unit]
Description=Thinkled

[Service]
Type=oneshot
ExecStart=thinkled
