# thinkled.service
[Unit]
Description=Thinkled

[Service]
Type=simple
ExecStart=/usr/local/sbin/thinkled

[Install]
WantedBy=multi-user.target
