# thinkled.service
[Unit]
Description=Thinkled (service)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/thinkled

[Install]
WantedBy=multi-user.target sleep.target
