# thinkled.service
[Unit]
Description=Thinkled (service)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/thinkled
