# thinkled/thinkled.service
[Unit]
Description=ThinkLed

[Service]
Type=oneshot
ExecStart=thinkled
