# beacon/tlbeacon.service
[Unit]
Description=ThinkLux Beacon

[Service]
Type=oneshot
ExecStart=tlbeacon
